# crontab
* * * * * /root/gather_data.sh
1 0,12 * * * /root/make_graphs.sh

